# cloudflared.service
[Unit]
Description=Cloudflare Tunnel for MailQueryWithMCP
After=network.target

[Service]
Type=simple
User=kimghw
Group=kimghw
WorkingDirectory=/home/kimghw/MailQueryWithMCP
ExecStart=/home/kimghw/.local/bin/cloudflared tunnel --config /home/kimghw/MailQueryWithMCP/cloudflare-tunnel-config.yml run
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudflared-mailquery

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/kimghw/.cloudflared /home/kimghw/MailQueryWithMCP/logs

[Install]
WantedBy=multi-user.target